# Turtlebot2/devel/etc/catkin/profile.d/25.turtlebot.sh
/workspaces/fyp/Turtlebot2/devel/.private/turtlebot_bringup/etc/catkin/profile.d/25.turtlebot.sh